# python/nav/etc/cron.d/navdhcp
## info: Collect graphite metrics from netboxes hosting a DHCP statistics server
*/1 * * * * navdhcp
